# .docker/config/supervisord.conf
[supervisord]
nodaemon=true

[program:dnsmasq]
command=dnsmasq --no-daemon --log-facility=/dev/stdout --conf-dir=/etc/dnsmasq.d,\*.conf
autorestart=true
startretries=999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm7.4
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:default-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /api/artisan queue:work --queue=default
numprocs=3
autostart=true
autorestart=true
user=rpi
stopwaitsecs=3600
startretries=999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:operations-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /api/artisan queue:work --queue=operations
numprocs=1
autostart=true
autorestart=true
user=rpi
stopwaitsecs=3600
startretries=999999999
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:schedule]
command=bash /schedule.sh
autorestart=true
user=rpi
startretries=999999999
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
